# supervisord.conf
[supervisord]
nodaemon=true
autorestart=true
logfile = /var/log/supervisord.log
logfile_maxbytes = 100MB

[program:cron]
autorestart=true
autorestart=true
command=cron -f