# config/supervisor.conf
[program:app]
command = uwsgi --ini /config/app.ini
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:nginx]
command = service nginx restart
autostart=true
autorestart=true